# crontab
0 7-17 * * * curl 'https://vytastock-iicnq.ondigitalocean.app'

# An empty line is required at the end of this file for a valid cron file.
